# deploy/neofi-heartbeat.service
[Unit]
Description=NeoFi Cloudflare Heartbeat
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/neofi/heartbeat.env
ExecStart=/usr/local/bin/neofi-heartbeat.sh
